# repo-mirror/rpm-repo-mirror.container
[Container]
ContainerName=rpm-repo-mirror
Image=localhost/rpm-repo-mirror
Volume=/data:/data:Z

[Service]
Restart=always